# meta-unisoc/recipes-update/swupdate/swupdate/default/swupdate-unisoc.service
[Unit]
Description=Unisoc SWUpdate daemon

[Service]
ExecStart=/bin/sh /usr/bin/update-init.sh

[Install]
WantedBy=multi-user.target
