# meta-unisoc/recipes-graphics/qtdemo/qtdemo.service
[Unit]
Description=qt demo test
RequiresMountsFor=/run

[Service]
User=root
ExecStart=/usr/bin/qtdemo-init.sh

[Install]
WantedBy=multi-user.target

