# meta-unisoc/recipes-devtools/charge/files/sl8521e-3in1/charge.service
[Unit]
Description=Unisoc NativeMMI Service

[Service]
ExecStart=/usr/bin/charge

[Install]
WantedBy=multi-user.target
Alias=charge.service
