# meta-unisoc/recipes-connectivity/netcon/files/udx710-module-op-pi/script/netcon.service
[Unit]
Description=netcon Daemon Service
After=sprdrild.target

[Service]
ExecStart=/usr/bin/netcon


Restart=always

[Install]
WantedBy=multi-user.target

