# meta-unisoc/recipes-core/ram/ram/ram.service
[Unit]
Description=ram service

[Service]
ExecStart= /bin/sh /usr/bin/ram-init

[Install]
WantedBy=multi-user.target
Alias=ram.service
