# meta-unisoc/recipes-connectivity/netcon/files/defaluts/script/netcon.service
[Unit]
Description=netcon Daemon Service
After=sprdrild.target

[Service]
ExecStart=/usr/bin/netcond


Restart=always

[Install]
WantedBy=multi-user.target

