# meta-unisoc/recipes-navigation/gnss-ota/files/gnss-ota.service
[Unit]
Description=gnss-ota

[Service]
ExecStart=/usr/bin/gnss-ota

[Install]
WantedBy=multi-user.target
Alias=gnss-ota.service
