# meta-unisoc/recipes-core/powerd/files/powerd.service
[Unit]
Description=powerd Service

[Service]
ExecStart=/usr/bin/powerd

Restart=always

[Install]
WantedBy=multi-user.target
Alias=powerd.service

