# meta-unisoc/recipes-devtools/factorytest/files/s9863a-1h10-otadev/factorytest.service
[Unit]
Description=Unisoc NativeMMI Service

[Service]
ExecStart=/usr/bin/factorytest

[Install]
WantedBy=multi-user.target
Alias=factorytest.service
