# meta-unisoc/recipes-security/tee_rpc/files/teerpc.service
[Unit]
Description=teerpc Service

[Service]
ExecStart=/usr/bin/teerpc

[Install]
WantedBy=multi-user.target
Alias=teerpc.service
